# docs/dev/tddy-daemon.service.example
# Example systemd unit for tddy-daemon.
# Runs unprivileged (User=tddy) with cgroup v2 delegation + an AppArmor userns grant, so the Linux
# sandbox works without root. Copy to /etc/systemd/system/tddy-daemon.service and adjust.
#
# Prefer automated install from the repo root: sudo ./install --systemd
# (optional: ./install --systemd --build to run ./release first). ./install also creates the tddy
# service user, chowns the log/auth dirs, and installs+loads the AppArmor profile. See
# daemon.yaml.production and docs/ft/daemon/systemd-install.md.
#
# Manual install: sudo cp tddy-daemon.service.example /etc/systemd/system/tddy-daemon.service
# Enable:  sudo systemctl enable tddy-daemon
# Start:   sudo systemctl start tddy-daemon

[Unit]
Description=tddy-daemon multi-user daemon
After=network.target

[Service]
Type=simple
ExecStart=/path/to/tddy-daemon -c /etc/tddy/daemon.yaml
Restart=on-failure
RestartSec=5

# Unprivileged operation: Delegate=yes gives the service a writable cgroup v2 subtree for the
# sandbox's per-session scopes; AppArmorProfile grants the unprivileged user namespaces the sandbox
# needs (load the profile first: apparmor_parser -r /etc/apparmor.d/tddy-daemon).
User=tddy
Group=tddy
Delegate=yes
AppArmorProfile=tddy-daemon
# NOTE: as a non-root user the daemon spawns sessions as itself (same_user=true). Multi-user setuid
# spawning requires root — set User=root/Group=root instead (Delegate=/AppArmorProfile= then optional).

[Install]
WantedBy=multi-user.target
